# ros_catkin_ws/devel/etc/catkin/profile.d/05.catkin_make.bash
/home/mars/MarquetteRMC2021/ros_catkin_ws/devel/.private/catkin/etc/catkin/profile.d/05.catkin_make.bash