# ros_catkin_ws/devel/etc/catkin/profile.d/1.ros_python_version.sh
/home/mars/MarquetteRMC2021/ros_catkin_ws/devel/.private/ros_environment/etc/catkin/profile.d/1.ros_python_version.sh