# ros_catkin_ws/devel/etc/catkin/profile.d/10.rosbuild.sh
# generated from rosbuild/env-hooks/10.rosbuild.sh.em

# env variables in develspace
export ROS_ROOT="/home/mars/MarquetteRMC2021/ros_catkin_ws/src/ros/rosbuild"
