# ros_catkin_ws/devel/etc/catkin/profile.d/10.rosconsole.sh
# generated from rosconsole/env-hooks/10.rosconsole.sh.develspace.em

export ROSCONSOLE_CONFIG_FILE="/home/mars/MarquetteRMC2021/ros_catkin_ws/src/rosconsole/config/rosconsole.config"
