# ros_catkin_ws/devel/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/mars/MarquetteRMC2021/ros_catkin_ws/src/ros/rosbash/rostcsh"
