# ros_catkin_ws/devel/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/home/mars/MarquetteRMC2021/ros_catkin_ws/src/ros/rosbash/roszsh"
