# ros_catkin_ws/devel/.private/ros_environment/etc/catkin/profile.d/1.ros_etc_dir.sh
# generated from ros_environment/env-hooks/1.ros_etc_dir.sh.em

# env variable in develspace
export ROS_ETC_DIR="/home/mars/MarquetteRMC2021/ros_catkin_ws/devel/.private/ros_environment/etc/ros"
